# btc-agent/deploy/systemd/agent.service
[Unit]
Description=BTC Agent V2
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=agent
Group=agent
WorkingDirectory=/var/lib/agent
EnvironmentFile=/etc/agent/agent.env
ExecStart=/opt/agent/current/agent scheduler --config /etc/agent/config.yaml
Restart=on-failure
RestartSec=10
TimeoutStopSec=45
KillSignal=SIGTERM
LimitNOFILE=65536
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/agent /var/log/agent
UMask=0077

[Install]
WantedBy=multi-user.target
