# btc-agent/deploy/systemd/btc-agent-web-console.service
[Unit]
Description=BTC Agent read-only Web Console
After=network-online.target btc-agent-immutable.service
Wants=network-online.target

[Service]
Type=simple
Environment=BTC_AGENT_WEB_DB_PATH=%h/btc-agent/runtime/data/btc_agent.db
Environment=BTC_AGENT_WEB_STATIC_DIR=%h/btc-agent/web-console/current/static
Environment=BTC_AGENT_WEB_LISTEN_ADDR=127.0.0.1:8787
ExecStart=%h/btc-agent/web-console/current/web-console
Restart=on-failure
RestartSec=10
TimeoutStopSec=30
KillSignal=SIGTERM
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only

[Install]
WantedBy=default.target
