# btc-agent/deploy/systemd/btc-agent-v2.service
[Unit]
Description=BTC Agent V2 authoritative scheduler
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/home/admin/btc-agent/btc-agent
EnvironmentFile=/home/admin/btc-agent-systemd.env
Environment=BTC_AGENT_INSTANCE_ID=v2-prod-01
Environment=BTC_AGENT_SCHEDULER_LOCK_FILE=/home/admin/.btc-agent-v2-scheduler.lock
ExecStart=/home/admin/btc-agent/btc-agent/bin/btc-agent scheduler --config /home/admin/btc-agent/btc-agent/config.yaml
Restart=on-failure
RestartSec=15
TimeoutStopSec=60
KillSignal=SIGTERM
LimitNOFILE=65536
UMask=0077
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
