# btc-agent/deploy/systemd/btc-agent-circuit-research.timer
[Unit]
Description=Run deterministic Circuit research every four hours

[Timer]
OnCalendar=*-*-* 00/4:17:00 UTC
RandomizedDelaySec=10m
Persistent=true
Unit=btc-agent-circuit-research.service

[Install]
WantedBy=timers.target
